# config/systemd/user/borg-backup.timer
[Unit]
Description=System backup with Borg

[Timer]
OnCalendar=daily UTC
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
